# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php]
command=%(ENV_SUPERVISOR_PHP_COMMAND)s
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; FIX Docker/Prioridad 4 (auditoría → corrección):
;
; Se verificó que routes/console.php define 3 tareas programadas
; (audit:verify-chain diario, purge-expired-dual-approvals diario,
; saldos:reconciliar --all diario), pero este contenedor no tenía cron ni
; ningún proceso ejecutando el scheduler de Laravel — las tareas estaban
; definidas pero JAMÁS se ejecutaban en este entorno Docker.
;
; El start-container de Sail invoca supervisord apuntando exactamente a
; /etc/supervisor/conf.d/supervisord.conf (un solo archivo, sin
; [include] hacia conf.d/*.conf), así que un archivo adicional en ese
; directorio no se carga solo — por eso este archivo REEMPLAZA (via bind
; mount en compose.yaml) al supervisord.conf original de la imagen de
; Sail, agregando el programa de abajo. El resto del archivo es idéntico
; al original (vendor/laravel/sail/runtimes/8.3/supervisord.conf) para no
; alterar el comportamiento existente del servidor PHP.
;
; `schedule:work` es el comando oficial de Laravel para correr el
; scheduler en un proceso long-running (alternativa a cron en contenedores).
[program:laravel-schedule]
command=php /var/www/html/artisan schedule:work
user=sail
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
